# root/etc/crontabs/root
## Camera Schedule Shooting ##
*/10 * * * * /usr/bin/capture
##
##
##
## Heartbeat PING ##
##
## DO NOT CHANGE THIS SETTING ##
*/3 * * * * /usr/bin/heartbeat
## DO NOT CHANGE THIS SETTING ##
##
