# root/etc/crontabs/root
## Camera Schedule Shooting ##
*/10 * * * * /usr/bin/capture
##
##
## Heartbeat PING ##
*/10 * * * * /usr/bin/heartbeat
##
##
