# root/etc/crontabs/root
## Heartbeat PING ##
##
## DO NOT CHANGE THIS SETTING ##
*/3 * * * * /usr/bin/heartbeat
## DO NOT CHANGE THIS SETTING ##
##
