# root/etc/uci-defaults/webcam
config general 'general'
    option latitude 'prova'
    option longitude 'prova'
    option faraday 'prova'
		
config shutter 'shutter'
    option imageformat 'prova'
		option imagesize 'prova'
		option imagequality 'prova'
		option whitebalance 'prova'

config day 'day'
    option iso_day 'prova'
		option aperture_day 'prova'
		option shutterspeed_day 'prova'

config night 'night'
    option iso_night 'prova'
		option aperture_night 'prova'
		option shutterspeed_night 'prova'

config server 'server'
    option upload_server 'prova'
		option upload_directory 'prova'
		option upload_username 'prova'
		option upload_password 'prova'
