# src/root/etc/crontabs/root
## Camera Schedule Shooting ##
*/10 * * * * /usr/bin/capture
##
##
## Heartbeat PING ##
## Do not change this setting ##
*/3 * * * * /usr/bin/heartbeat
## Do not change this setting ##
##
##
